# rootfs/etc/s6/crond/run
#!/bin/sh

exec /usr/bin/doas /usr/sbin/crond -f -L /dev/stdout
